# deploy/clash-tracker-ingestion.timer
[Unit]
Description=Run Clash Tracker ingestion every 10 minutes

[Timer]
OnCalendar=*-*-* *:0/10:00
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
